# libexec/supervisord/boards.conf
[program:boards]
environment=PORT="57518"
command=envdir ../../etc/teamstatus.d/ node core/app.js
directory=libexec/boards
redirect_stderr=true
stdout_logfile=logs/boards.log
