# etc/supervisor.d/boards.conf
[program:boards]
environment=PORT="57518"
command=foreman start web -f lib/boards/Procfile -e etc/teamstatus.conf
stdout_logfile=logs/boards.log