# etc/supervisor.d/boards.conf
[program:boards]
environment=PORT="57518"
command=foreman start -f lib/boards/Procfile -e etc/teamstatus.conf
redirect_stderr=true
stdout_logfile=logs/boards.log