# etc/supervisor.d/console.conf
[program:console]
environment=PORT="57517"
command=foreman start -f lib/console/Procfile -e etc/teamstatus.conf
redirect_stderr=true
stdout_logfile=logs/console.log