# etc/supervisor.d/console.conf
[program:console]
command=envdir ../../etc/teamstatus.d/ bundle exec unicorn -o 127.0.0.1 -p 57517 -c ./config/unicorn.rb
stdout_logfile=logs/console.log
directory=lib/console